# deployment/water-tank-control.service
[Unit]
Description=Hot Water Tank Temperature Control System
After=network.target evok.service
Wants=evok.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/water-tank-control
Environment="PATH=/opt/water-tank-control/venv/bin"
Environment="USE_MOCK_EVOK=false"
ExecStart=/opt/water-tank-control/venv/bin/python /opt/water-tank-control/app.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=water-tank-control

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/water-tank-control/logs /opt/water-tank-control

[Install]
WantedBy=multi-user.target
