# deployment/water-tank-control.service
[Unit]
Description=Hot Water Tank Temperature Control System
After=network.target evok.service
Wants=evok.service

[Service]
Type=simple
# User and Group are automatically set by install.sh to match the installing user
# If installing manually, change these to match your system user
User=jarek
Group=jarek
WorkingDirectory=/opt/water-tank-control
Environment="PATH=/opt/water-tank-control/venv/bin"
Environment="USE_MOCK_EVOK=false"
# Use gunicorn with eventlet worker for production (supports WebSocket)
# For standalone mode (no nginx): --bind 0.0.0.0:5000
# For nginx reverse proxy: --bind 127.0.0.1:5000
# Note: Eventlet worker required for WebSocket support with Flask-SocketIO
ExecStart=/opt/water-tank-control/venv/bin/gunicorn --worker-class eventlet -w 1 --bind 127.0.0.1:5000 app:app
# Alternative if eventlet has issues: use gevent worker
# ExecStart=/opt/water-tank-control/venv/bin/gunicorn --worker-class gevent -w 1 --bind 127.0.0.1:5000 app:app
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=water-tank-control

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/water-tank-control/logs /opt/water-tank-control

[Install]
WantedBy=multi-user.target
